# services/watchdog_microsite.service
[Unit]
Description=Monitor microsite and restart it if offline
After=network.target

[Service]
ExecStart=/usr/local/bin/watchdog.sh
Restart=always
RestartSec=300

# Use the dedicated microsite system user
User=microsite
Group=microsite

# Allow USB control via uhubctl
ReadWritePaths=/dev/bus/usb

# Restrict modifications to system files
ProtectSystem=full
ProtectHome=yes

# Direct service logs to systemd journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
